# scripts/motion-capture-chick.service
[Unit]
Description=Motion-capture timelapse (Chick Cam)
After=camera-hls-3.service

[Service]
Type=simple
User=ajsornig
ExecStart=/bin/bash /home/ajsornig/chicken-stream/scripts/motion-capture.sh chick
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
